# etc/hotplug.d/iface/40-igmpproxy
#!/bin/sh
[ "$INTERFACE" = wan ] && [ "$ACTION" = ifup ] && /etc/init.d/igmpproxy restart
